# tools/timesfm_server/timesfm-api.service.example
[Unit]
Description=TimesFM forecasting API for aibig9
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/min_hs/min_hs/minun001.github.io/tools/timesfm_server
Environment=TIMESFM_API_TOKEN=CHANGE_ME
Environment=TIMESFM_ALLOWED_ORIGINS=https://minun001.github.io,http://127.0.0.1:4000,http://localhost:4000
Environment=TIMESFM_MODEL_ID=google/timesfm-2.5-200m-pytorch
Environment=TIMESFM_DEVICE=auto
Environment=HF_HOME=/data/min_hs/hf_cache/timesfm
Environment=MAX_UPLOAD_MB=25
ExecStart=/home/min_hs/min_hs/timesfm-api/.venv/bin/uvicorn app:app --host 127.0.0.1 --port 8791
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
